# src/lib/games/ToonTownRewritten/supervisord.conf
[supervisord]
nodaemon=true

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1280x720x24 -ac
stdout_logfile=/tmp/xvfb.log
priority=10
autorestart=true
redirect_stderr=true

[program:x11vnc]
environment=DISPLAY=:99
command=/usr/bin/x11vnc -display :99 -rfbauth /root/.vnc/passwd -forever -shared -noxdamage -repeat -rfbport 5900 -o /tmp/x11vnc.log
stdout_logfile=/tmp/x11vnc_stdout.log
# start after openbox and the app are up so VNC captures actual windows
priority=35
autorestart=true
redirect_stderr=true

[program:openbox]
environment=DISPLAY=:99
command=/usr/bin/openbox-session
stdout_logfile=/tmp/openbox.log
priority=20
autorestart=true
redirect_stderr=true

[program:ttr]
environment=DISPLAY=:99
command=/usr/bin/xterm -hold -e /launcher/ttr_launch.sh
stdout_logfile=/tmp/ttr.log
priority=30
autorestart=true
redirect_stderr=true

